# rootfs/etc/services.d/firefox/run
#!/usr/bin/with-contenv sh

su $PROFILE_NAME -c "/usr/bin/firefox -P \"$PROFILE_NAME\" --new-tab --url $WEBSITE"
